# deploy/supervisor/mtgstore-worker.conf
; Supervisor program for the Symfony Messenger worker (alternative to systemd).
; Copy to /etc/supervisor/conf.d/, edit paths/user/environment, then:
;   supervisorctl reread && supervisorctl update
;
; numprocs runs several workers in parallel (they claim rows with SKIP LOCKED,
; so parallelism is safe). Each worker self-exits on the time/memory limit and
; supervisor restarts it, bounding long-lived PHP memory growth.

[program:mtgstore-worker]
command=/usr/bin/php bin/console messenger:consume async failed --time-limit=3600 --memory-limit=256M --failure-limit=10 -v
directory=/var/www/mtgstore/backend
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopsignal=TERM
; APP_ENV=prod plus the injected secrets (keep real secrets out of git — point
; this at your secrets file or set them in the supervisor environment).
environment=APP_ENV="prod"
stdout_logfile=/var/log/mtgstore/worker.log
stderr_logfile=/var/log/mtgstore/worker.err.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
